# dotfiles/.config/systemd/user/files.service
[Unit]
Description=Run checkfiles program 

[Service]
Type=simple
ExecStart=%h/.local/bin/rofi_launcher --reindex

[Install]
WantedBy=default.target
